# .config/profile.d/alacritty.sh
#!/bin/sh

#
# ~/.config/sh/profile.d/hx.sh
# Alacritty configurations login shell script
#

export TERMINAL=alacritty
