# .config/profile.d/clang.sh
#!/bin/sh

#
# ~/.config/sh/profile.d/clang.sh
# Clang configurations login shell script
#

export CC=clang
export CXX=clang++
