# .config/profile.d/hx.sh
#!/bin/sh

#
# ~/.config/sh/profile.d/hx.sh
# Helix configurations login shell script
#

export VISUAL=hx
