# .config/profile.d/sqlite.sh
#!/bin/sh

#
# ~/.config/sh/profile.d/sqlite.sh
# SQLite configurations login shell script
#

export SQLITE_HISTORY="${XDG_STATE_HOME:-$HOME/.local/state}"/sqlite_history
