# .config/sh/profile.d/clang.sh
#!/bin/sh

#
# ~/.config/sh/profile.d/clang.sh
# Clang configurations login shell script
#

# TODO: Not really a clang configuration.
export CC=clang
export CXX=clang++
